# logs.service
[Unit]
Description=Run the logs service

[Service]
ExecStart=/usr/local/bin/logs.sh


